# deploy/supervisord/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/autoguard/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ============================================================
; Redis Server
; ============================================================
[program:redis]
command=/usr/bin/redis-server --port 6379 --daemonize no --maxmemory 256mb --maxmemory-policy allkeys-lru
user=autoguard
autostart=true
autorestart=true
priority=100
stdout_logfile=/var/log/autoguard/redis.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/autoguard/redis-error.log
stderr_logfile_maxbytes=10MB

; ============================================================
; Nginx Server
; ============================================================
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=200
stdout_logfile=/var/log/nginx/access.log
stderr_logfile=/var/log/nginx/error.log

; ============================================================
; Next.js Web Dashboard (Port 3000)
; ============================================================
[program:web]
command=node /app/packages/web/.next/standalone/server.js
directory=/app/packages/web
user=autoguard
environment=NODE_ENV="production",PORT="3000",HOSTNAME="127.0.0.1"
autostart=true
autorestart=true
priority=300
startretries=5
startsecs=10
stdout_logfile=/var/log/autoguard/web.log
stdout_logfile_maxbytes=50MB
stderr_logfile=/var/log/autoguard/web-error.log
stderr_logfile_maxbytes=50MB

; ============================================================
; Cloak Worker (Port 3001)
; ============================================================
[program:cloak-worker]
command=node /app/workers/cloak-worker/dist/index.js
directory=/app/workers/cloak-worker
user=autoguard
environment=NODE_ENV="production",PORT="3001"
autostart=true
autorestart=true
priority=300
startretries=5
startsecs=10
stdout_logfile=/var/log/autoguard/cloak-worker.log
stdout_logfile_maxbytes=50MB
stderr_logfile=/var/log/autoguard/cloak-worker-error.log
stderr_logfile_maxbytes=50MB

; ============================================================
; Log Writer (Background Process)
; ============================================================
[program:log-writer]
command=node /app/workers/log-writer/dist/index.js
directory=/app/workers/log-writer
user=autoguard
environment=NODE_ENV="production"
autostart=true
autorestart=true
priority=400
startretries=5
startsecs=5
stdout_logfile=/var/log/autoguard/log-writer.log
stdout_logfile_maxbytes=50MB
stderr_logfile=/var/log/autoguard/log-writer-error.log
stderr_logfile_maxbytes=50MB

; ============================================================
; Scheduler Worker (Background Process)
; Handles: blacklist sync, expiry cleanup, stats aggregation
; ============================================================
[program:scheduler]
command=node /app/workers/scheduler/dist/index.js
directory=/app/workers/scheduler
user=autoguard
environment=NODE_ENV="production",BLACKLIST_SYNC_INTERVAL="300000",EXPIRY_CLEANUP_INTERVAL="3600000",STATS_AGGREGATION_INTERVAL="300000"
autostart=true
autorestart=true
priority=500
startretries=5
startsecs=10
stdout_logfile=/var/log/autoguard/scheduler.log
stdout_logfile_maxbytes=50MB
stderr_logfile=/var/log/autoguard/scheduler-error.log
stderr_logfile_maxbytes=50MB

; ============================================================
; Queue Worker (Background Process)
; Handles: page generation, scraping, AI generation tasks
; ============================================================
[program:queue-worker]
command=node /app/workers/queue-worker/dist/index.js
directory=/app/workers/queue-worker
user=autoguard
environment=NODE_ENV="production",QUEUE_POLL_INTERVAL="1000",MAX_CONCURRENT_JOBS="2"
autostart=true
autorestart=true
priority=500
startretries=5
startsecs=10
stdout_logfile=/var/log/autoguard/queue-worker.log
stdout_logfile_maxbytes=50MB
stderr_logfile=/var/log/autoguard/queue-worker-error.log
stderr_logfile_maxbytes=50MB

; ============================================================
; Process Groups
; ============================================================
[group:autoguard]
programs=redis,nginx,web,cloak-worker,log-writer,scheduler,queue-worker
priority=999
